# systemd/stow/.config/user/polkit-gnome.service
[Unit]
Description=PolicyKit Authentication Agent
PartOf=graphical-session.target

[Service]
ExecStart=/usr/lib/polkit-gnome-authentication-agent-1
Restart=on-failure

[Install]
WantedBy=graphical-session.target
